# systemd/zwetow-check-update.timer
[Unit]
Description=Periodic Zwetow update check

[Timer]
OnBootSec=2min
OnUnitActiveSec=30min
Persistent=true

[Install]
WantedBy=timers.target
