# systemd/zwetow-check-update.service
[Unit]
Description=Check for Zwetow updates

[Service]
Type=oneshot
ExecStart=/opt/zwetow/bin/zwetow-check-update.sh
